# monroe-experiments/etc/cron.d/monroe-experiments
# /etc/cron.d/monroe-experiments
# sync monroe experiments, ensuring they are running

*/1 * * * *   root    /usr/bin/monroe-experiments

