# monroe-experiments/etc/cron.d/monroe-experiments
# /etc/cron.d/monroe-experiments
# sync monroe experiments, ensuring they are running

*/1 * * * *   root  /usr/bin/monroe-experiments
27 */4 * * *  root  /usr/bin/pull-base-containers
@reboot       root  /bin/rm /var/lock/monroe-experiments.lock || true 

