# sysevent/DEBIAN/postinst
#!/bin/sh

systemctl enable sysevent-shutdown || true
systemctl start  sysevent-shutdown || true
