# sysevent/DEBIAN/prerm
#!/bin/sh

systemctl stop    sysevent-shutdown || true
systemctl disable sysevent-shutdown || true

