# tuptime/DEBIAN/postinst
#!/bin/sh

systemctl enable tuptime || true
systemctl start  tuptime || true
