# tuptime/DEBIAN/prerm
#!/bin/sh

systemctl stop    tuptime || true
systemctl disable tuptime || true

