# usage/etc/cron.d/usage
# /etc/cron.d/usage
# record data usage on relevant interfaces and namespaces

*/1 * * * *   root    /usr/bin/usage /monroe/usage '' eth0 wlan0 wwan0 usb0 usb1 usb2
*/1 * * * *   root    /usr/bin/usage ip netns exec monroe /monroe/usage/netns /tmp/metadata-usage-monroe eth0 op0 op1 op2 
