# usage/etc/cron.d/usage
# /etc/cron.d/usage
# record data usage on relevant interfaces and namespaces

*/1 * * * *   root    /usr/bin/usage-defaults
