# ubusd/deb_dist/DEBIAN/postinst
#!/bin/sh

systemctl enable ubusd || true
systemctl start  ubusd || true
