# ubusd/deb_dist/DEBIAN/prerm
#!/bin/sh

systemctl stop ubusd 
systemctl disable ubusd
