# mf910-switcher/etc/cron.d/mf910-switcher
# /etc/cron.d/mf910-switcher
# modeswitch MF910 modems

*/1 * * * *   root    /usr/bin/mf910-switcher

