# metadata-exporter/DEBIAN/prerm
#!/bin/sh

systemctl stop metadata-timeloop 
systemctl disable metadata-timeloop 

systemctl stop metadata-exporter
systemctl disable metadata-exporter 
