# wifimetadata-publisher/DEBIAN/postinst
#!/bin/sh

systemctl enable wifimetadata-publisher || true
systemctl start  wifimetadata-publisher || true
