# wifimetadata-publisher/DEBIAN/prerm
#!/bin/sh

systemctl stop wifimetadata-publisher 
systemctl disable wifimetadata-publisher