# metadata-exporter/lib/systemd/system/metadata-exporter.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=Metadata Exporter
After=network.target

[Service]
EnvironmentFile=/etc/default/metadata-exporter

ExecStart=/usr/sbin/metadata-exporter -c /etc/metadata-exporter.conf

Type=simple
Restart=on-failure

StandardError=null
StandardOutput=null
